# files/systemd/bucketlist-api-gcp.service
[Unit]
Description=bucketlist-api-gcp
[Service]
ExecStart=/var/www/bucketlist-api-gcp/start.sh

[Install]
WantedBy=multi-user.target