# files/supervisor/bucketlist-api-gcp.conf
[program:bucketlist-api-gcp]
command = /var/www/bucketlist-api-gcp/start.sh
autostart=true
autorestart=trueq